# ops/swarming-api.service
# The network's public join/submit API — a SEPARATE process from
# swarming-board.service (the World Cup showcase). Same box, same domain,
# same database, different port: a bug or load spike here can never affect
# the WC show's process. Install:
#   sudo cp ops/swarming-api.service /etc/systemd/system/
#   sudo mkdir -p /etc/systemd/system/swarming-api.service.d
#   sudo cp ops/swarming-api.override.conf /etc/systemd/system/swarming-api.service.d/override.conf
#   sudo systemctl daemon-reload && sudo systemctl enable --now swarming-api
[Unit]
Description=Swarming public network API (isolated from the WC showcase)
After=network-online.target

[Service]
User=ec2-user
WorkingDirectory=/opt/swarming
# Deliberately NOT /etc/swarming/env (that file pins SWARMING_PORT=8090 for
# the WC show and would win over an Environment= override here in practice,
# despite unit-file ordering suggesting otherwise). This service needs no
# provider secrets, so it gets its own fully explicit environment instead.
Environment=SWARMING_DB=/opt/swarming/server/data/swarming-live.db
Environment=SWARMING_PORT=8091
Environment=SWARMING_HOST=127.0.0.1
Environment=SWARMING_TRUST_PROXY=1
ExecStart=/usr/bin/node server/src/index.ts
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
